# templates/folder.service.erb
[Unit]
Description=Bittorent Sync service (<%= @path %>)
After=network.target

[Service]
Type=simple
User=<%= @owner %>
Group=<%= @group %>
UMask=<%= @umask %>
ExecStart=/usr/bin/btsync --nodaemon --config <%= @config %>
TimeoutStopSec=10
Restart=on-abort

[Install]
WantedBy=multi-user.target
